# debian/conf/power-status-changed.conf
start on power-status-changed

task

exec shutdown -h now "Power Down"
