# extra/conf-session/upstart-dconf-bridge.conf
description	"Bridge dconf events into session upstart"

emits dconf

start on started dbus
stop on stopped dbus

respawn

exec upstart-dconf-bridge
