# debian/upstart-bin.upstart.cron.daily
#!/bin/sh

# For each Upstart Session Init, emit "rotate-logs" event, requesting
# the session Inits to rotate their logs. There is no user-daily cron.
#
# Doing it this way does not rely on System Upstart, nor
# upstart-event-bridge(8) running in the Session Init.
#
# Note that system-level Upstart logs are handled separately using a
# logrotate script.

[ -x /sbin/initctl ] || exit 0

for file in /run/user/*/upstart/sessions/*.session
do
    [ -f "$file" ] || continue
    [ -h "$file" ] && continue

    session=$(grep \
	    "^UPSTART_SESSION=unix:abstract=/com/ubuntu/upstart-session/[0-9][0-9]*/[0-9][0-9]*$" \
	    "$file" 2>/dev/null || true)
    [ -z "$session" ] && continue

    env -i "$session" /sbin/initctl emit rotate-logs >/dev/null 2>&1 || true
done
